# lab3/configs/prod/aise.service
[Unit]
Description=aise service with docker-compose.
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=always
TimeoutSec=300

WorkingDirectory=/opt/aise

# Compose up
ExecStart=docker compose -f docker-compose.yml up

# Compose down, remove containers
ExecStop=docker compose -f docker-compose.yml down

# Compose restart
ExecReload=docker compose -f docker-compose.yml restart


[Install]
WantedBy=multi-user.target